# supervisord.conf
[supervisord]
nodaemon=true


[program:pureftpd]
command=/usr/sbin/vsftpd
autorestart=true

[program:nodeweb]
command=node app
directory=/app
